# packaging/debian/globus-gram-job-manager-pbs/debian/globus-gram-job-manager-pbs-setup-seg.postinst
#!/bin/sh

set -e

if [ -x "/etc/init.d/globus-scheduler-event-generator" ]; then
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d globus-scheduler-event-generator start pbs || :
	else
		/etc/init.d/globus-scheduler-event-generator start pbs || :
	fi
fi

#DEBHELPER#
