# packaging/debian/globus-gatekeeper/debian/globus-gatekeeper.prerm
#!/bin/sh
set -e

if [ -x "/etc/init.d/globus-gatekeeper" ]; then
	invoke-rc.d globus-gatekeeper stop || exit 1
fi

#DEBHELPER#
