# packaging/debian/globus-scheduler-event-generator/debian/globus-scheduler-event-generator-progs.prerm
#!/bin/sh
set -e

if [ -x "/etc/init.d/globus-scheduler-event-generator" ]; then
	invoke-rc.d globus-scheduler-event-generator stop || exit 1
fi

#DEBHELPER#
