# packaging/debian/globus-scheduler-event-generator/debian/globus-scheduler-event-generator-progs.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -x "/etc/init.d/globus-scheduler-event-generator" ]; then
		update-rc.d globus-scheduler-event-generator defaults >/dev/null
		invoke-rc.d globus-scheduler-event-generator start || exit 1
	fi
fi

#DEBHELPER#
