# packaging/debian/globus-gridftp-server/debian/globus-gridftp-server-progs.prerm
#!/bin/sh
set -e

if [ -x "/etc/init.d/globus-gridftp-server" ]; then
	invoke-rc.d globus-gridftp-server stop || exit 1
	invoke-rc.d globus-gridftp-sshftp stop || exit 1
fi

#DEBHELPER#
